# crontab
0 2 * * 0 /Users/vincent/developer/macos-config/update-optimized.sh >> /Users/vincent/.local/share/update-optimized.log 2>&1
